# samples/launchd/ai.itsup.bringup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- macOS equivalent of itsup-bringup.service. Loads the guardian wrapper,
     which runs `itsup run && itsup apply` then waits on a signal. launchd
     sends SIGTERM at logout/shutdown; the guardian traps it and runs
     `itsup down --clean` for graceful container teardown — same end-state
     as systemd's ExecStop, via launchd's actual semantics.

     KeepAlive[SuccessfulExit]=false: a clean SIGTERM exit (shutdown) is NOT
     a reason to relaunch; only a crash is. ExitTimeOut=120 gives the
     cleanup enough time to tear down all stacks before launchd SIGKILLs. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.itsup.bringup</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ROOT}}/bin/bringup-guardian.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ExitTimeOut</key>
    <integer>120</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <!-- launchd's default PATH omits Homebrew/Docker locations; the guardian
             invokes itsup, which shells out to docker/sops/etc via PATH, so set
             one covering Apple Silicon + Intel Homebrew. -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/itsup.bringup.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/itsup.bringup.log</string>
</dict>
</plist>
