# samples/launchd/ai.itsup.apply.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- macOS equivalent of itsup-apply.timer + itsup-apply.service. Fires nightly
     at 03:00 local time and runs `itsup apply`. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.itsup.apply</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ROOT}}/.venv/bin/itsup</string>
        <string>apply</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>ITSUP_ROOT</key>
        <string>{{ROOT}}</string>
        <!-- launchd's default PATH omits Homebrew/Docker locations; itsup shells
             out to docker/sops/etc via PATH, so set one covering Apple Silicon +
             Intel Homebrew. -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/itsup.apply.log</string>
    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/itsup.apply.log</string>
</dict>
</plist>
