# samples/launchd/ai.itsup.healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- macOS equivalent of pi-healthcheck.timer + pi-healthcheck.service.
     Fires every 5 minutes and runs the healthcheck script (which also
     writes to /var/log/instrukt-ai/itsup/pi-healthcheck.log via tee). -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.itsup.healthcheck</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ROOT}}/bin/pi-healthcheck.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{HOME}}</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/var/log/instrukt-ai/itsup/pi-healthcheck.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/instrukt-ai/itsup/pi-healthcheck.log</string>
</dict>
</plist>
