# samples/systemd/itsup-backup.service
[Unit]
Description=itsUP nightly backup to S3 (runs bin/backup.py)
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
# Runs as root so backup can read every container's bind-mounted volume
# regardless of which UID inside the container wrote the files (e.g.
# AdGuard writes its conf as root in-container, which appears root-owned
# on the host — morriz can't read it).
User=root
WorkingDirectory={{ROOT}}
ExecStart=/bin/bash -lc 'cd {{ROOT}} && source env.sh && .venv/bin/python bin/backup.py'
