# samples/systemd/itsup-bringup.service
[Unit]
Description=Orchestrate full itsUP stack lifecycle (start on boot, stop on shutdown)
After=docker.service network-online.target
Requires=docker.service network-online.target
PartOf=docker.service

[Service]
Type=oneshot
User={{USER}}
Group={{GROUP}}
WorkingDirectory={{ROOT}}
Environment=HOME=/home/{{USER}}
Environment=ITSUP_ROOT={{ROOT}}
# Start core infra then deploy all projects (apply may take time for pulls).
# Type=oneshot runs ExecStartPost only after ExecStart succeeds — the systemd-native
# form of the old `run && apply` (systemd does not interpret `&&`).
ExecStart={{ROOT}}/.venv/bin/itsup run
ExecStartPost={{ROOT}}/.venv/bin/itsup apply
# On shutdown, stop containers in dependency order but DO NOT remove them —
# leaves them in `exited` state so Docker's restart policy brings them back
# sub-second on next boot (vs minutes to recreate from compose).
ExecStop={{ROOT}}/.venv/bin/itsup down
TimeoutStopSec=120
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
